# aws/root/etc/init/update_ssmtp_fqdn.conf
# Keep route53 updated
#

description     "Update /etc/ssmtp/ssmtp.conf"

console output

start on net-device-up

task
exec /usr/local/bin/update_ssmtp_fqdn.sh
