# aws/root/etc/init/update_route53_cname.conf
# Keep route53 updated
#

description     "Update route53 cname"

console output

start on net-device-up

task
exec /usr/local/bin/update_route53_cname.sh
