# aws/root/etc/cron.daily/tomcat-cleanup
#!/bin/sh

if [ -d /usr/local/tomcat/logs ] ; then
/usr/bin/find /usr/local/tomcat/logs/ -type f -mtime +30 -delete 2>/dev/null
fi

