# aws/servers/portal-generic/update_wordpress_hostname.conf
# Keep route53 updated
#

description     "Update /etc/wordpress/config-master.php"

console output

start on net-device-up

task
exec /usr/local/bin/update_wordpress_hostname.sh
